# services/bpffs
#!/sbin/openrc-run

# this is an excerpt taken from the sysfs script in openrc edge
# see https://github.com/OpenRC/openrc/blob/77a49de9cd72c669ef136985cdfaa1570b7ac896/init.d/sysfs.in#L109 

description="Mount the bpf filesystem."
sysfs_opts=nodev,noexec,nosuid
start()
{
    # set up kernel support for bpf filesystem
	if [ -d /sys/fs/bpf ] && ! mountinfo -q /sys/fs/bpf; then
		if grep -qs 'bpf$' /proc/filesystems; then
			ebegin "Mounting bpf filesystem"
			mount -n -t bpf -o ${sysfs_opts} bpf /sys/fs/bpf
			eend $?
		fi
	fi
    return 0
}
