# services/bridge_setup
#!/sbin/openrc-run

description="Establish a host-VM bridge."
sysfs_opts=nodev,noexec,nosuid
start()
{
    ebegin $description

    sysctl -w net.ipv4.ip_forward=1

    ip addr add 192.168.101.100/24 dev eth0
    ip link set eth0 up

    ip addr add 192.168.102.100/24 dev eth1
    ip link set eth1 up

    eend $?

    return 0
}
